# service/zenstorage.service
[Unit]
Description=ZenStorage Mount
After=network-online.target

[Service]
User=seed
Group=seed
Type=notify
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/rclone mount \
  --config=/home/seed/.config/rclone/rclone.conf \
  --allow-other \
  --allow-non-empty \
  --rc \
  --rc-addr=localhost:5576 \
  --vfs-read-ahead=128M \
  --vfs-read-chunk-size=64M \
  --vfs-read-chunk-size-limit=2G \
  --vfs-cache-mode=full \
  --vfs-cache-max-age=24h \
  --vfs-cache-max-size=100G \
  --fast-list \
  --buffer-size=64M \
  --dir-cache-time=1h \
  --timeout=10m \
  --tpslimit=10 \
  --tpslimit-burst=100 \
  --umask=002 \
  --async-read=false \
  --syslog \
  -v \
  zenstorage: /mnt/zenstorage
ExecStop=/bin/fusermount -uz /mnt/zenstorage
Restart=on-abort
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=default.target
