# service/zenstorage_primer.timer
[Unit]
Description=Zenstorage Primer - Timer

[Timer]
OnUnitInactiveSec=8760h

[Install]
WantedBy=timers.target
