# service/zenstorage_primer.timer
[Unit]
Description=Zenstorage Primer - Timer

[Timer]
OnUnitInactiveSec=167h

[Install]
WantedBy=timers.target
